# content/service/3/run
#!/bin/sh

exec 2>&1 >/dev/null

exec /frp/frps -c /frp/frps.toml
